# server/supervisor/development.ini
[program:omop_harvest_development]
directory = /home/devuser/webapps/omop_harvest_env/omop_harvest
command = /home/devuser/webapps/omop_harvest_env/bin/uwsgi --ini /home/devuser/webapps/omop_harvest_env/omop_harvest/server/uwsgi/development.ini
autostart = true
autorestart = true
