# server/supervisor/development.ini
[program:omop_harvest-development]
directory = /home/devuser/webapps/omop_harvest-env/omop_harvest
command = /home/devuser/webapps/omop_harvest-env/bin/uwsgi --ini /home/devuser/webapps/omop_harvest-env/omop_harvest/server/uwsgi/development.ini
autostart = true
autorestart = true
