# server/supervisor/development.ini
[program:cbttc-harvest-development]
directory = /home/devuser/webapps/cbttc-harvest-env/cbttc-harvest
command = /home/devuser/webapps/cbttc-harvest-env/bin/uwsgi --ini /home/devuser/webapps/cbttc-harvest-env/cbttc-harvest/server/uwsgi/development.ini
autostart = true
autorestart = true
