# server/supervisor/local.ini
[program:cbttc-local]
directory = /home/devuser/.virtualenvs/cbttc-env/cbttc
command = /home/devuser/.virtualenvs/cbttc-env/bin/uwsgi --ini /home/devuser/.virtualenvs/cbttc-env/cbttc/server/uwsgi/local.ini
autostart = true
autorestart = true
