# server/supervisor/local.ini
[program:omop_harvest-local]
directory = /home/devuser/.virtualenvs/omop_harvest-env/omop_harvest
command = /home/devuser/.virtualenvs/omop_harvest-env/bin/uwsgi --ini /home/devuser/.virtualenvs/omop_harvest-env/omop_harvest/server/uwsgi/local.ini
autostart = true
autorestart = true
