# server/supervisor/local.ini
[program:omop_harvest]
directory = /Users/brownea/Projects/omop_harvest
command = /Users/brownea/.virtualenvs/omop_harvest/bin/uwsgi --ini /Users/brownea/Projects/omop_harvest/server/uwsgi/local.ini
autostart = true
autorestart = true
