# server/supervisor/production.ini
[program:omop_harvest-production]
directory = /home/devuser/webapps/omop_harvest-env/omop_harvest
command = /home/devuser/webapps/omop_harvest-env/bin/uwsgi --ini /home/devuser/webapps/omop_harvest-env/omop_harvest/server/uwsgi/production.ini
autostart = true
autorestart = true
