# server/supervisor/production.ini
[program:cbttc-harvest-production]
directory = /home/brpuser/webapps/cbttc-harvest-env/cbttc-harvest
command = /usr/local/bin/uwsgi --ini /home/brpuser/webapps/cbttc-harvest-env/cbttc-harvest/server/uwsgi/production.ini
autostart = true
autorestart = true
