# server/supervisor/production.ini
[program:omop_harvest_production]
directory = /home/brpuser/webapps/omop_harvest_env/omop_harvest
command = /usr/local/bin/uwsgi --ini /home/brpuser/webapps/omop_harvest_env/omop_harvest/server/uwsgi/production.ini
autostart = true
autorestart = true
