# server/supervisor/vagrant.ini
[program:cbttc-harvest]
directory = /home/vagrant/webapps/cbttc-env/cbttc-harvest
command = /home/vagrant/webapps/cbttc-env/bin/uwsgi --ini /home/vagrant/webapps/cbttc-env/cbttc-harvest/server/uwsgi/vagrant.ini
autostart = true
autorestart = true
