# cid/lib/logstash/supervisor_logstash_run.conf
[supervisord]
nodaemon=false

[program:logstash_omop_harvest]
directory=/opt/logstash
startretries=3
redirect_stderr = true
user=root
command=/opt/logstash/bin/logstash -f /opt/apps/harvest-app/cid/research.chop.edu/logstash/omop_harvest.logstash.conf
autostart=true
startsecs=5
